# com.claude.dashboard-frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude.dashboard-frontend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npx</string>
    <string>vite</string>
    <string>--port</string>
    <string>4000</string>
    <string>--strictPort</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/lucaalice/Desktop/ClaudexSEO/agent-dashboard/frontend</string>
  <key>StandardOutPath</key>
  <string>/tmp/claude-dashboard-frontend.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claude-dashboard-frontend.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
